# payload/systemd/ubuntu-zombie-health.timer
[Unit]
Description=Run Ubuntu Zombie health check every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=ubuntu-zombie-health.service

[Install]
WantedBy=timers.target
