# payload/systemd/ubuntu-zombie-health.service
[Unit]
Description=Periodic Ubuntu Zombie health check

[Service]
Type=oneshot
User=__AGENT_USER__
Group=__AGENT_USER__
ExecStart=/opt/ai-zombie/bin/health-check
